# 5.6-nginx/config/supervisord.conf
[supervisord]
nodaemon=true

[program:php-fpm]
command=bin/bash -c "exec entrypoint php-fpm"

[program:nginx]
command=/bin/bash -c "exec entrypoint nginx"
